# packaging/pilothoused.service
[Unit]
Description=Pilothouse privileged broker
After=systemd-sysext.service
Before=pilothouse.service

[Service]
Type=simple
ExecStart=/usr/local/libexec/pilothoused --socket /run/pilothouse/broker.sock --socket-group pilothouse --admin-group sudo
Restart=on-failure
RestartSec=3
User=root
Group=pilothouse
RuntimeDirectory=pilothouse
RuntimeDirectoryMode=0750
UMask=0077
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes

[Install]
WantedBy=multi-user.target
